# lib/imaging/layouts/common/etc/systemd/system/chrony.service.d/override.conf
# SPDX-License-Identifier: Apache-2.0
# Copyright 2020 Kirill 'kkm' Katsnelson
#
# This file was installed by BurrMill.
#
# chronyd distribution does not install this provider dependency on the
# time-sync.target, rendering it unavailable. We do not have dependecies on it,
# by in case we ever do, this is a fix.

[Unit]
Before=time-sync.target
Wants=time-sync.target
