# lib/imaging/layouts/common/etc/systemd/system/systemd-remount-fs.service.d/override.conf
# This file was installed by BurrMill.

# Workaround for https://github.com/systemd/systemd/issues/14603 (fix #14618).
# remount-fs and growfs@ race so that grows may fail because the root filesystem
# is still mounted readonly. This ensures the remount happens first.
[Unit]
Before=systemd-growfs@-.service
