# lib/imaging/layouts/common/etc/systemd/system/ssh.service.d/override.conf
# This file was installed by BurrMill.
#
# Make sure host keys are available before sshd can start. If that fails, do not
# restart immediately, hold off a bit.

[Service]
ExecStartPre=
ExecStartPre=/usr/local/sbin/update_host_keys $UHKEY_VERBOSE
ExecStartPre=/usr/sbin/sshd -t
RestartSec=3
